# init-scripts/ubuntu1404/node_exporter.conf
# Run node_exporter
start on startup

script
   /usr/local/bin/node_exporter
end script
